# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=docker-php-entrypoint php-fpm
autostart=true
autorestart=true
priority=5

[program:queue]
command=php artisan queue:work --tries=3
autostart=true
autorestart=true
priority=10
